# awsggpi4.service
[Unit]
Description=YOLO Object Detection Stream (awsggpi4)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/awsggpi4
EnvironmentFile=-/etc/default/awsggpi4
ExecStart=%h/awsggpi4/bin/python %h/awsggpi4/objectdetection.py
Restart=on-failure
RestartSec=3
StandardOutput=append:%h/awsggpi4/awsggpi4.log
StandardError=append:%h/awsggpi4/awsggpi4-error.log

[Install]
WantedBy=multi-user.target
